# backend/freight-registry.service
[Unit]
Description=Freight Equasis registry crawler (vessel data enrichment)
After=network.target

[Service]
Type=oneshot
User=lbzgiu
WorkingDirectory=/home/lbzgiu/quant/freight/backend
ExecStart=/home/lbzgiu/quant/freight/backend/.venv/bin/python -m registry.crawl
EnvironmentFile=-/home/lbzgiu/quant/freight/backend/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
