# backend/freight-registry.timer
[Unit]
Description=Run Equasis registry crawler once daily
Requires=freight-registry.service

[Timer]
OnBootSec=15min
OnCalendar=*-*-* 04:30:00
Persistent=true

[Install]
WantedBy=timers.target
